# build/crontab
0 0 * * SUN /root/update-unifi.sh > /var/log/update-unifi.log
